# etc/cron.d/awscli-tools
# THIS FILE IS CONTROLLED BY GIT.  DO NOT EDIT BY HAND OR YOUR CHANGES
# WILL DISAPPEAR.
#
# Run the deploy once per day
25 1 * * * root /foresight/bin/cronmail -u -a foresight-se@fsautomation.com -- /usr/local/sbin/awscli-tools-run-deploy

# Update the identity files every hour
44 * * * * root /foresight/bin/cronmail -u -a foresight-se@fsautomation.com -- /usr/local/sbin/awscli-tools-make-identity-files
